# overlay/composer/etc/cron.monthly/composer
#!/usr/bin/env bash

/usr/local/bin/composer self-update
