# config/clamav/etc/cron.daily/clamav
#!/bin/sh

freshclam && clamscan -ir /var/www/ --log=/var/log/clamscan.log --move=/var/quarantine/ --scan-mail=no
