# overlay/iptables/etc/network/if-pre-up.d/iptables
#!/usr/bin/env bash

iptables-restore < /etc/iptables.rules && exit 0
